# .config/systemd/user/adapt_backlight.service
[Unit]
Description=Adaptive Backlight Daemon

[Service]
ExecStart=/home/d0np0br3/Scripts/adapt_backlight.sh

[Install]
WantedBy=default.target
